# addons_updater_enhanced/rootfs/etc/services.d/addons-updater/run
#!/usr/bin/with-contenv bash
exec /usr/bin/addons-updater 2>&1
