# zerotier-controller/s6/zerotier/run
#!/usr/bin/with-contenv bash
exec /run.sh